# docker/elk/v1.1/etc/supervisor/conf.d/beats.conf
[program:beats]
command=/bin/bash load.sh
