# docker/elk/v1.1/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=service ngnix restart
